# ops/com.pcpartscrawler.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pcpartscrawler.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/joejaeyoung/2026/pc-parts-crawler/.venv/bin/python</string>
    <string>scripts/run_all.py</string>
    <string>--skip-sources</string>
    <string>daangn</string>
    <string>--danawa-pages</string>
    <string>0</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/joejaeyoung/2026/pc-parts-crawler</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
    <key>CRAWL_TRIGGER_SOURCE</key>
    <string>launchd</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/joejaeyoung/Library/Logs/pc-parts-crawler/daily.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/joejaeyoung/Library/Logs/pc-parts-crawler/daily.err</string>
</dict>
</plist>
